# services/sshd/run
#!/usr/bin/with-contenv bash
# SSH_PASSWORD 未设置则不启动
if [ -z "$SSH_PASSWORD" ]; then
    echo "[sshd] SSH_PASSWORD not set, SSH disabled"
    exec sleep infinity
fi
# 持久化 host keys → /config/.ssh-host-keys/
KEY_DIR="/config/.ssh-host-keys"
mkdir -p "$KEY_DIR"
[ -f "$KEY_DIR/ssh_host_rsa_key" ] || ssh-keygen -t rsa -b 4096 -f "$KEY_DIR/ssh_host_rsa_key" -N "" -q
[ -f "$KEY_DIR/ssh_host_ed25519_key" ] || ssh-keygen -t ed25519 -f "$KEY_DIR/ssh_host_ed25519_key" -N "" -q
chown root:root "$KEY_DIR"/*_key; chmod 600 "$KEY_DIR"/*_key
# 设置 abc 密码 + 确保有 login shell
echo "abc:${SSH_PASSWORD}" | chpasswd
usermod -s /bin/bash abc 2>/dev/null || true
echo "[sshd] Starting (port 22)..."
exec /usr/sbin/sshd -D -e
