# services/docker-ready/run
#!/usr/bin/with-contenv bash
# s6 服务：Docker 就绪检测

if [ "$ENABLE_DIND" != "true" ]; then
    echo "[docker-ready] DinD disabled, skipping..."
    sleep infinity
fi

if ! command -v docker &> /dev/null; then
    echo "[docker-ready] Docker not installed, skipping..."
    sleep infinity
fi

log() { echo "[docker-ready] $1"; }
log "Waiting for Docker daemon..."

for i in {1..30}; do
    if docker info > /dev/null 2>&1; then
        log "Docker daemon is ready!"
        docker version --format 'Docker: {{.Server.Version}}' 2>/dev/null || true
        while docker info > /dev/null 2>&1; do
            sleep 10
        done
        log "Docker daemon stopped!"
        exit 1
    fi
    sleep 1
done

log "Timeout waiting for Docker"
exit 1
